# systemd/user/udiskie.service
[Unit]
Description=udiskie - Automount removable media
PartOf=graphical-session.target
Requires=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/udiskie --no-config
Restart=on-failure
RestartSec=5
Slice=session.slice
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
